# service/deluge-web/run
#!/bin/sh
exec /sbin/setuser seedbox /usr/bin/deluge-web -d -c /deluge/deluge-web -p 8080 --loglevel=info -l /deluge/deluge-web/webui.log
