# service/deluge-web/run
#!/bin/sh
exec /sbin/setuser deluge /usr/bin/deluge-web -c /deluge/deluge-web -p 8080 --loglevel=info -l /deluge/deluge-web/webui.log
