# service/deluged/run
#!/bin/sh
exec /sbin/setuser seedbox /usr/bin/deluged -d -c /deluge/deluged --loglevel=info -l /deluge/deluged/deluged.log
